# smart_monitor.service
[Unit]
Description=Smart Security Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/smart_monitor/smart_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
